# Repository of Kaeon/2 - Collection/1 - Implementation/1 - APIs/4 - Utilities/2 - Hardware/2 - Kaeon Control/2 - GHI/1 - Source/2 - Modules/1 - Adapters/1 - Raspberry Pi/2 - Scheduling/crontab
@reboot sudo /usr/local/bin/node /home/pi/Documents/GHI/ghiSource.js 1233